# macos/launchd/com.bjanda.obsidian-new-year.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bjanda.obsidian-new-year</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__HOME__/.local/bin/obsidian-new-year.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Month</key>
    <integer>1</integer>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/obsidian-new-year.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/obsidian-new-year.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
