# macos/launchd/local.stats-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.stats-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/.local/bin/stats-push.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/stats-push.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/stats-push.log</string>
</dict>
</plist>
