# linux/systemd/backup-proxmox.service
[Unit]
Description=Proxmox config backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-proxmox.sh
