# linux/systemd/backup-pi.service
[Unit]
Description=Pi server backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-pi.sh
