# linux/systemd/backup-arch.service
[Unit]
Description=Arch-server backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-arch.sh
