# linux/systemd/freshrss-autoupdate.service
[Unit]
Description=Auto-update whitelisted services via FreshRSS feeds

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freshrss-autoupdate.sh
