# linux/systemd/freshrss-maintenance.service
[Unit]
Description=Nightly FreshRSS & RSS-Bridge maintenance

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freshrss-maintenance.sh
