# linux/systemd/freshrss-autoupdate.timer
[Unit]
Description=Daily auto-update check for whitelisted services

[Timer]
OnCalendar=*-*-* 09:00
Persistent=true

[Install]
WantedBy=timers.target
