# linux/systemd/gpx-manifest.path
[Unit]
Description=Watch for GPX file changes

[Path]
PathChanged=/var/www/hiking-map/gpx

[Install]
WantedBy=multi-user.target
