# linux/systemd/gpx-manifest.service
[Unit]
Description=Watch for GPX file changes and rebuild manifest
After=local-fs.target

[Service]
Type=simple
User=admin
ExecStart=/usr/local/bin/gpx-manifest-watch.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
