# linux/systemd/health-check-proxmox.service
[Unit]
Description=Proxmox health check with ntfy alerts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health-check-proxmox.sh
