# linux/systemd/health-check.service
[Unit]
Description=Pi health check with ntfy alerts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pi-health-check.sh
