# linux/systemd/health-check-arch.timer
[Unit]
Description=Daily arch-server health check

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
