# linux/systemd/health-check.service
[Unit]
Description=Health check with ntfy alerts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health-check.sh
