# linux/systemd/health-check-pi.timer
[Unit]
Description=Daily Pi health check

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
