# linux/systemd/health-check-pi.service
[Unit]
Description=Pi health check with ntfy alerts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health-check-pi.sh
