# linux/systemd/health-check-arch.service
[Unit]
Description=Arch-server health check with ntfy alerts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health-check-arch.sh
