# linux/systemd/ntfy-failure@.service
[Unit]
Description=Send ntfy alert on %i failure

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s -o /dev/null \
    -H "Title: %i failed" \
    -H "Priority: high" \
    -H "Tags: skull,warning" \
    -d "Service %i failed on pi-server" \
    http://localhost:2586/pi-alerts
