# linux/systemd/pi-backup.service
[Unit]
Description=Pi server backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pi-backup.sh
