# linux/systemd/pkg-maintenance.service
[Unit]
Description=Weekly package update, upgrade, and cleanup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/home/admin
ExecStart=/home/admin/.local/bin/pkg-maintenance.sh
