# linux/systemd/pkg-maintenance.timer
[Unit]
Description=Run package maintenance weekly on Sunday

[Timer]
OnCalendar=Sun 09:00
Persistent=true

[Install]
WantedBy=timers.target
