# linux/systemd/qbt-cleanup.timer
[Unit]
Description=Monthly qBittorrent seeding cleanup

[Timer]
OnCalendar=*-*-01 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
