# linux/systemd/rss-bridge-cache-cleanup.timer
[Unit]
Description=Daily RSS-Bridge cache cleanup

[Timer]
OnCalendar=*-*-* 04:00
Persistent=true

[Install]
WantedBy=timers.target
