# linux/systemd/tailscale-certrenew.service
[Unit]
Description=Renew Tailscale HTTPS certificates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tailscale-certrenew.sh
