# linux/systemd/wifi-watchdog.timer
[Unit]
Description=Run WiFi watchdog every 2 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
